# scripts/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command = /usr/sbin/sshd -D

[program:build]
command = cd /opt/app/eums/eums/client git&& grunt build-staging:${LOCAL_HOST_IP}

[program:postgresql]
command= /opt/scripts/startPostgres.sh
autostart = true

[program:redis]
command = redis-server
autostart = true

[program:mongo]
command = mongod
autostart = true

[program:contacts]
command = /opt/scripts/startContacts.sh
autostart = true

[program:uwsgi]
command = /usr/local/bin/uwsgi --emperor /etc/uwsgi/sites
autostart = true

[program:nginx]
command = service nginx start
autostart = true
